# etc/s6-overlay/s6-rc.d/controller/run
#!/bin/sh
exec s6-setuidgid oasis /usr/local/bin/controller
